# priv/ansible/roles/redis/files/upstart_redis.conf
start on runlevel [2345]
stop on runlevel [!2345]
console log
respawn
setuid redis
exec /usr/local/bin/redis-server --bind 127.0.0.1 --dir /var/lib/redis